# deploy/serverhub.supervisor.conf
; Supervisor program for the ServerHub panel itself.
; {PANEL_ROOT} and {PANEL_USER} are substituted by deploy/install.sh.
[program:serverhub]
command={PANEL_ROOT}/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8765
directory={PANEL_ROOT}/backend
user={PANEL_USER}
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/serverhub.err.log
stdout_logfile=/var/log/supervisor/serverhub.out.log
